# deployment/sem-platform.service
[Unit]
Description=SEM Segmentation Statistics Platform - Backend API
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/root/wd_bishe
ExecStartPre=/usr/bin/docker compose -f /root/wd_bishe/deployment/docker-compose.yml pull
ExecStart=/usr/bin/docker compose -f /root/wd_bishe/deployment/docker-compose.yml up -d --build
ExecStop=/usr/bin/docker compose -f /root/wd_bishe/deployment/docker-compose.yml down
ExecReload=/usr/bin/docker compose -f /root/wd_bishe/deployment/docker-compose.yml up -d --build
TimeoutStartSec=0
Restart=on-failure
RestartSec=30s

# 环境变量文件
EnvironmentFile=/root/wd_bishe/deployment/.env

[Install]
WantedBy=multi-user.target
